# dracut/30ignition/coreos-populate-var.service
[Unit]
Description=Populate OSTree /var
DefaultDependencies=false

# Need to do this with all mount points active
After=ignition-mount.service

# But *before* we start dumping files in there
Before=ignition-files.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/coreos-populate-var
